# debian/opi-backend.service
[Unit]
Description=Backend daemon for Keep

[Service]
ExecStart=/usr/sbin/opi-backend -d
Type=forking

[Install]
WantedBy=multi-user.target
