# debian/opi-backend.service
[Unit]
Description=Backend daemon for Keep

[Service]
ExecStart=/usr/sbin/opi-backend -d
PIDFile=/var/run/opi-backend.pid
Type=forking

[Install]
WantedBy=multi-user.target
